# podkop/files/etc/hotplug.d/iface/50-podkop
#!/bin/sh

/etc/init.d/podkop add_route_interface
